# mount.conf
# Ubuntu upstart file at /etc/init/mount.conf

start on vagrant-mounted
 
script
  mount --bind /home/vagrant/app/.meteor /vagrant/app/.meteor
end script